# ai_labs.service
[Unit]
Description=AI Labs Admin Panel - Gunicorn Application
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/AI_Labs_Portal
Environment="PATH=/var/www/AI_Labs_Portal/venv/bin"
ExecStart=/var/www/AI_Labs_Portal/venv/bin/gunicorn --workers 3 --bind 0.0.0.0:8000 'app:create_app()'

[Install]
WantedBy=multi-user.target
